# roles/admin_bastion/templates/bastion.timer.j2
[Unit]
Description=Check the bastion configuration regularly

[Timer]
OnBootSec=15min
OnUnitActiveSec=2h
Unit=bastion.service

[Install]
WantedBy=multi-user.target